# gui/mmdvm-DVconfig.service
[Unit]
Description=MMDVM Status Service
After=network-online.target mariadb.service
Wants=network-online.target mariadb.service
Requires=mariadb.service

[Service]
Type=simple
User=mmdvm
Group=mmdvm
WorkingDirectory=/usr/local/bin
# warte bis der Socket existiert (max ~30s)
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do [ -S /run/mysqld/mysqld.sock ] && exit 0; sleep 1; done; echo "mysqld.sock fehlt"; exit 1'
ExecStart=/usr/local/bin/DVconfig
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
